# launchd/com.nathan.claude-obsidian-capture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nathan.claude-obsidian-capture</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/nathan/Code/Atelier/Projects/claude-obsidian-capture/scripts/run-capture.sh</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>WorkingDirectory</key>
  <string>/Users/nathan/Code/Atelier/Projects/claude-obsidian-capture</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Seoul</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/nathan/Code/Atelier/Projects/claude-obsidian-capture/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nathan/Code/Atelier/Projects/claude-obsidian-capture/logs/launchd.err.log</string>
</dict>
</plist>
